# image/mkosi.extra/etc/systemd/system/icecast2.service.d/override.conf
[Service]
ExecStartPre=/usr/local/bin/icecast-setup.sh
